# mbp_fan_control.service
[Unit]
Description=Controls Fan speed on Macbook 466

[Service]
Type=simple
ExecStart=/usr/local/bin/mbp_fan_control
